# crontab
5 8 * * 1-5 root /usr/local/bin/radikorec --duration 224 --prefix SUPPIN --channel NHK1 --directory /radikorec
5 8 * * 6 root /usr/local/bin/radikorec --duration 40 --prefix BUNGEI --channel NHK1 --directory /radikorec

