# crontab
5 8 * * 1-5 root /usr/local/bin/rec_radiru.sh R1 224 /radikorec SUPPIN
5 8 * * 6 root /usr/local/bin/rec_radiru.sh R1 40 /radikorec BUNGEI

